# etc/supervisor/conf.d/pymbserver.conf
[program:pymbserver]
command=pymbserver -H '' -p 5020
user=pi
autostart=true
autorestart=false
environment=PYTHONUNBUFFERED=1
