# etc/supervisor/conf.d/pymbserver.conf
[program:pymbserver]
command=pymbserver -H 0.0.0.0 -p 5020
user=pi
autostart=true
autorestart=false
environment=PYTHONUNBUFFERED=1
